# files/supervisor/getconf.conf
[program:vcpe-getconf]
command=/usr/local/bin/getconf -ip core
autostart=true
startsecs=5
user=root
