# files/supervisor/ipsecping.conf
[program:vcpe-ipsecping]
command=/usr/local/bin/ipsecPing
autostart=true
startsecs=5
user=root
