# files/supervisor/policyupdate.conf
[program:vcpe-cpepolicyupdate]
command=/usr/local/bin/policyUpdate
autostart=true
startsecs=5
user=root
