# files/supervisor/policylatency.conf
[program:vcpe-cpe-policy-latecy]
command=/usr/local/bin/policylatency
autostart=true
startsecs=5
user=root
