# files/supervisor/vcpe-cperecovery.conf
[program:vcpe-cperecovery]
command=/usr/local/bin/cpeRecovery
autostart=true
startsecs=10
user=root
