# files/supervisor/vcpe-ipsec.conf
[program:vcpe-ipsec]
command=/usr/sbin/ipsec start --nofork
autostart=true
startsecs=10
user=root
