# files/supervisor/vcpe-polipo.conf
[program:vcpe-polipo]
command=/usr/local/bin/polipo -c /etc/polipo.conf
autostart=true
startsecs=10
user=root
