# files/supervisor/vcpe-frpc.conf
[program:vcpe-frpc]
command=/usr/local/bin/frpc -c /etc/frpc.ini
autostart=true
startsecs=10
user=root
