# files/daho_vcpe_recovery
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=17
STOP=97

boot() {
 if [ -f /usr/bin/cpe_Recovery ] ; then
      /usr/bin/cpe_Recovery
 fi
}
